# src/StartMesure.service
[Unit]
Description = hoge
After=local-fs.target
ConditionPathExists = /opt/sakulab_ex_modlue/src

[Service]
ExecStart=/opt/sakulab_ex_module/src/start_program.sh
Restart=always
Type=simple

[Install]
WantedBy=multi-user.target